# src/controller/service/controller/run
#!/bin/sh
echo "Starting controller..."
exec ping -c 8 google.com
