# src/controller/service/app/run
#!/bin/sh
echo "Starting app..."
exec ping -c 8 google.com
